# s6/debian-root/etc/services.d/cloudflared/finish
#!/usr/bin/with-contenv bash

s6-echo "Stopping cloudflared"
kill -9 $(pgrep cloudflared)
